# configs/systemd/cmux-ios-vnc-proxy.service
[Unit]
Description=Cmux iOS simulator VNC websocket proxy (reverse mode)
After=network.target

[Service]
Type=simple
Environment=CMUX_VNC_PROXY_PORT=39387
Environment=CMUX_VNC_TARGET_PORT=39386
Environment=CMUX_VNC_REVERSE=true
Environment=CMUX_VNC_WEB_DIR=/usr/share/novnc
ExecStartPre=/bin/mkdir -p /var/log/cmux
ExecStart=/usr/local/lib/cmux/cmux-vnc-proxy
Restart=always
RestartSec=3
StandardOutput=append:/var/log/cmux/ios-vnc-proxy.log
StandardError=append:/var/log/cmux/ios-vnc-proxy.log

[Install]
WantedBy=cmux.target
